# systemd/pi-vpn-gateway.service
[Unit]
Description=Pi VPN Gateway — AP routing orchestration
After=network-online.target hostapd.service openvpn-brazil.service
Wants=network-online.target openvpn-brazil.service
Requires=hostapd.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/pi-vpn-gateway/env
ExecStart=/opt/pi-vpn-gateway/systemd/gateway-start.sh
ExecStop=/opt/pi-vpn-gateway/systemd/gateway-stop.sh

[Install]
WantedBy=multi-user.target
